# crontab
*/5 * * * * $HOME/bin/email >> $HOME/.cron_email_log
