# crontab
*/5 * * * * /home/urbain/bin/email >> /home/urbain/.cron_email_log
