# crontab
*/10 * * * * /usr/bin/offlineimap -u quiet
