# crontab
*/10 * * * * /usr/bin/offlineimap -u quiet
*/5 * * * * sh ~/.xmodmap
