# crontab/crontab
*/5 * * * * DISPLAY=:0 /usr/bin/env mbsync -a
*/5 * * * * /usr/bin/env mr push
*/5 * * * * /usr/bin/env mr update
0 0 * * * /usr/bin/pacman -Syuw --noconfirm
