# crontab/crontab
*/5 * * * * DISPLAY=:0 /usr/bin/env mbsync -a
0 11 * * * /usr/bin/env mr push
0 11 * * * /usr/bin/env mr update
0 11 * * * /usr/bin/pacman -Syuw --noconfirm
