# crontab/crontab
*/5 * * * * /home/urbain/bin/email >> /home/urbain/.cron_email_log
*/5 * * * * /usr/bin/vdirsyncer sync >> /home/urbain/.cron_calendar_log
