# crontab/crontab
*/5 * * * * DISPLAY=:0 /usr/bin/env mbsync -a
